# wicd/wicd.service
[Unit]
Description=Wicd Network Manager

[Service]
Type=dbus
BusName=org.wicd.daemon
ExecStart=/usr/sbin/wicd -f

[Install]
WantedBy=network.target
